# webserver/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command=php-fpm -D

[program:apache]
command=httpd -D FOREGROUND